# cmd/tuna/oyster-init-server.conf
[program:oyster-init-server]
command=/home/ubuntu/oyster-init-server --vsock-addr 3:1300 --job-id placeholder_job_id
autostart=true
autorestart=true

